# 02.threat-stack/timers/blackice-summary.timer
[Unit]
Description=🕗 Set BlackIce Daily Summary to run at 07:00

[Timer]
OnCalendar=*-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
