# 02.threat-stack/services/blackice-summary.service
[Unit]
Description=BlackIce Daily Summary Report
Wants=blackice-summary.timer

[Service]
Type=oneshot
ExecStart=/usr/local/bin/blackice-daily-summary.sh
User=root
Group=root
